# server/caden-api.service
[Unit]
Description=Caden API — local AI gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/caden/server
ExecStart=/home/pi/caden/server/venv/bin/python main.py
Restart=always
RestartSec=5
EnvironmentFile=/home/pi/caden/server/.env

[Install]
WantedBy=multi-user.target
